# qsdk/target/linux/ipq806x/base-files/etc/uci-defaults/network
#!/bin/sh
#
# Copyright (c) 2015 The Linux Foundation. All rights reserved.
# Copyright (C) 2011 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh
. /lib/ipq806x.sh

ucidef_set_interface_loopback

board=$(ipq806x_board_name)

case "$board" in
db149|\
db149_1xx|\
db149_2xx)
	ucidef_set_interfaces_lan_wan "eth1 eth2 eth3" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "6 1 2 3 4"
	ucidef_add_switch_vlan "switch0" "2" "0 5"
	;;
ap148|\
ap148_1xx|\
ap145|\
ap145_1xx)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "6 1 2 3 4"
	ucidef_add_switch_vlan "switch0" "2" "0 5"
	;;
ap160|\
ap161)
	ucidef_set_interfaces_lan_wan "eth1 eth2" "eth0"
	ucidef_add_switch "switch0" "1" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "6 4"
	ucidef_add_switch_vlan "switch0" "2" "0 1 2 3"
	;;
ap160_2xx)
	ucidef_set_interface_bond "bond0" "2" "eth1 eth2 eth3" "2"
	ucidef_set_interface_lan "bond0 eth0"
	ucidef_add_switch "switch0" "1" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "1 26"
	ucidef_add_switch_vlan "switch0" "2" "2 3 4 27"
	;;
storm)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "6 1 2 3"
	ucidef_add_switch_vlan "switch0" "2" "0 5"
	;;
whirlwind)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "6 1"
	ucidef_add_switch_vlan "switch0" "2" "0 2"
	;;
ak01_1xx)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	;;
ap-dk01.1-c1|\
ap-dk01.1-c2|\
ap-dk04.1-c1|\
ap-dk04.1-c2|\
ap-dk04.1-c3|\
ap-dk04.1-c5|\
ap-dk06.1-c1|\
ap-dk07.1-c1|\
ap-dk07.1-c2)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0t 1 2 3 4"
	ucidef_add_switch_vlan "switch0" "2" "0t 5"
	ucidef_add_switch_ext "switch0" "QosPtMode" "1" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "2" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "3" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "4" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "5" "dscp" "enable"

uci batch <<EOF
set network.guest_lan='interface'
set network.guest_lan.ifname='ath01'
set network.guest_lan.force_link=1
set network.guest_lan.type='bridge'
set network.guest_lan.proto='static'
set network.guest_lan.ipaddr='192.168.2.1'
set network.guest_lan.netmask='255.255.255.0'
set network.guest_lan.ip6assign='60'

set dhcp.guest_lan="dhcp"
set dhcp.guest_lan.interface='guest_lan'
set dhcp.guest_lan.start='100'
set dhcp.guest_lan.limit='150'
set dhcp.guest_lan.leasetime='12h'
set dhcp.guest_lan.dhcpv6='server'
set dhcp.guest_lan.ra='server'

add firewall zone
set firewall.@zone[-1].name='guest_lan'
add_list firewall.@zone[-1].network='guest_lan'
set firewall.@zone[-1].input='ACCEPT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='ACCEPT'
add firewall forwarding
set firewall.@forwarding[-1].src='guest_lan'
set firewall.@forwarding[-1].dest='wan'

set network.subdev_lan='interface'
set network.subdev_lan.ifname='ath02'
set network.subdev_lan.force_link=1
set network.subdev_lan.type='bridge'
set network.subdev_lan.proto='static'
set network.subdev_lan.ipaddr='192.168.3.1'
set network.subdev_lan.netmask='255.255.255.0'
set network.subdev_lan.ip6assign='60'

set dhcp.subdev_lan="dhcp"
set dhcp.subdev_lan.interface='subdev_lan'
set dhcp.subdev_lan.start='10'
set dhcp.subdev_lan.limit='250'
set dhcp.subdev_lan.leasetime='24h'
set dhcp.subdev_lan.dhcpv6='server'
set dhcp.subdev_lan.ra='server'

add firewall zone
set firewall.@zone[-1].name='subdev_lan'
add_list firewall.@zone[-1].network='subdev_lan'
set firewall.@zone[-1].input='ACCEPT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='ACCEPT'
add firewall forwarding
set firewall.@forwarding[-1].src='subdev_lan'
set firewall.@forwarding[-1].dest='wan'

commit dhcp
commit firewall

EOF


	;;
ap-dk04.1-c4)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0t 1 2 3"
	ucidef_add_switch_vlan "switch0" "2" "0t 4 5"
	ucidef_add_switch_ext "switch0" "QosPtMode" "1" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "2" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "3" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "4" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "5" "dscp" "enable"
	ucidef_add_switch_ext_trunk "switch0" "TrunkGroup" "0" "enable" "0x30"
	;;
ap-dk05.1-c1)
	ucidef_set_interface_lan "eth0"
	ucidef_add_switch "switch0" "1" "0"
	ucidef_add_switch_vlan "switch0" "0" "0 1 2 3 4 5"
	ucidef_add_switch_ext "switch0" "QosPtMode" "1" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "2" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "3" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "4" "dscp" "enable"
	ucidef_add_switch_ext "switch0" "QosPtMode" "5" "dscp" "enable"
	;;
hk01)
	ucidef_set_interfaces_lan_wan "eth1 eth2" "eth0"
	;;
*)
	echo "Unsupported hardware. Network interfaces not intialized"
	;;
esac

uci commit network

exit 0
